# private_dot_config/systemd/user/wlsunset.service
[Unit]
Description=Run wlsunset on startup

[Service]
ExecStart=/usr/bin/wlsunset -S 05:30 -s 17:40
Type=simple

[Install]
WantedBy=default.target

